# recipes-arrow-sockit/arrow-sockit-xfce-init/files/arrow-sockit-xfce-init.service
[Unit]
Description=Arrow SoCKit XFCE Session
After=arrow-sockit-fpga-init.service

[Service]
Type=simple
ExecStart=/usr/bin/startxfce4

[Install]
WantedBy=basic.target
