# recipes-arrow-sockit/arrow-sockit-fpga-init/files/arrow-sockit-fpga-init.service
[Unit]
Description=FPGA Initialization

[Service]
Type=simple
ExecStart=/usr/bin/arrow-sockit-fpga-init.sh

[Install]
WantedBy=basic.target
